# wifi-package/wifi_shutdown_util/files/wifi_shutdown
#!/bin/sh /etc/rc.common

STOP=21
USE_PROCD=1

shutdown() {
	/etc/init.d/prplmesh shutdown
	/etc/init.d/prplmesh_whm shutdown
}
